# deploy/com.kelvin.lifeos.capture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kelvin.lifeos.capture</string>

  <!-- Run the Telegram capture daemon under the repo's venv. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/kelvintan/SynologyDrive/Documents/life-os/.venv/bin/python</string>
    <string>/Users/kelvintan/SynologyDrive/Documents/life-os/capture_daemon.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os</string>

  <!-- Restart if it ever exits (network blips, crashes). -->
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <!-- PATH MUST include ~/.local/bin (where `claude` lives) and the nvm node bin,
       or the agentic router's `claude -p` dies with FileNotFoundError under launchd
       and every message falls back to #unsorted (the exact bug that stranded a
       capture for >1h). ffmpeg is in /opt/homebrew/bin for voice notes.
       Backstop: router.py/claude_cli.py also resolve `claude` by absolute path. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/kelvintan/.local/bin:/Users/kelvintan/.nvm/versions/node/v20.20.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os/data/capture.daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os/data/capture.daemon.err.log</string>
</dict>
</plist>
