# deploy/com.kelvin.lifeos.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kelvin.lifeos.web</string>

  <!-- Run the Flask dashboard on port 5070 under the repo's venv. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/kelvintan/SynologyDrive/Documents/life-os/.venv/bin/python</string>
    <string>/Users/kelvintan/SynologyDrive/Documents/life-os/server.py</string>
    <string>--port</string>
    <string>5070</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os</string>

  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os/data/web.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kelvintan/SynologyDrive/Documents/life-os/data/web.err.log</string>
</dict>
</plist>
